# web_flask/gunicorn.conf
description "Gunicorn application server handling hanzi"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data
chdir /home/ubuntu/hanzi

exec gunicorn --workers 3 --bind 0.0.0.0:5000 --bind 0.0.0.0:5000 wsgi hanzi.wsgi:application